# docker/s6-rc.d/nextjs/run
#!/command/with-contenv sh
cd /app
exec s6-setuidgid node node server.js
